# autossh.conf
# autosshd
#

description	"autosshd: manage all autossh ressions"

start on started udev
stop on stopped udev

task

script
	autosshd-reload
end script


